# profile.d/cuda.sh
# libraries added via ldconfig instead of LD_LIBRARY_PATH
#  export LD_LIBRARY_PATH=/usr/local/cuda-9.0/lib64:$LD_LIBRARY_PATH
export PATH=/usr/local/cuda-9.0/bin/:$PATH

