# profile.d/dircolors.sh
LS_COLORS='rs=0'
LS_COLORS+=':di=01;34'
LS_COLORS+=':ln=01;36'
LS_COLORS+=':mh=00'
LS_COLORS+=':pi=33'
LS_COLORS+=':so=01;35'
LS_COLORS+=':do=01;35'
LS_COLORS+=':bd=33;01'
LS_COLORS+=':cd=33;01'
LS_COLORS+=':or=31;01'
LS_COLORS+=':su=37;41'
LS_COLORS+=':sg=30;43'
LS_COLORS+=':ca=30;41'
LS_COLORS+=':tw=30;42'
LS_COLORS+=':ow=34;42'
LS_COLORS+=':st=37;44'
LS_COLORS+=':ex=01;32'
LS_COLORS+=':*.tar=01;31'
LS_COLORS+=':*.tgz=01;31'
LS_COLORS+=':*.arc=01;31'
LS_COLORS+=':*.arj=01;31'
LS_COLORS+=':*.taz=01;31'
LS_COLORS+=':*.lha=01;31'
LS_COLORS+=':*.lz4=01;31'
LS_COLORS+=':*.lzh=01;31'
LS_COLORS+=':*.lzma=01;31'
LS_COLORS+=':*.tlz=01;31'
LS_COLORS+=':*.txz=01;31'
LS_COLORS+=':*.tzo=01;31'
LS_COLORS+=':*.t7z=01;31'
LS_COLORS+=':*.zip=01;31'
LS_COLORS+=':*.z=01;31'
LS_COLORS+=':*.Z=01;31'
LS_COLORS+=':*.dz=01;31'
LS_COLORS+=':*.gz=01;31'
LS_COLORS+=':*.lrz=01;31'
LS_COLORS+=':*.lz=01;31'
LS_COLORS+=':*.lzo=01;31'
LS_COLORS+=':*.xz=01;31'
LS_COLORS+=':*.bz2=01;31'
LS_COLORS+=':*.bz=01;31'
LS_COLORS+=':*.tbz=01;31'
LS_COLORS+=':*.tbz2=01;31'
LS_COLORS+=':*.tz=01;31'
LS_COLORS+=':*.deb=01;31'
LS_COLORS+=':*.rpm=01;31'
LS_COLORS+=':*.jar=01;31'
LS_COLORS+=':*.war=01;31'
LS_COLORS+=':*.ear=01;31'
LS_COLORS+=':*.sar=01;31'
LS_COLORS+=':*.rar=01;31'
LS_COLORS+=':*.alz=01;31'
LS_COLORS+=':*.ace=01;31'
LS_COLORS+=':*.zoo=01;31'
LS_COLORS+=':*.cpio=01;31'
LS_COLORS+=':*.7z=01;31'
LS_COLORS+=':*.rz=01;31'
LS_COLORS+=':*.cab=01;31'
LS_COLORS+=':*.jpg=01;35'
LS_COLORS+=':*.jpeg=01;35'
LS_COLORS+=':*.gif=01;35'
LS_COLORS+=':*.bmp=01;35'
LS_COLORS+=':*.pbm=01;35'
LS_COLORS+=':*.pgm=01;35'
LS_COLORS+=':*.ppm=01;35'
LS_COLORS+=':*.tga=01;35'
LS_COLORS+=':*.xbm=01;35'
LS_COLORS+=':*.xpm=01;35'
LS_COLORS+=':*.tif=01;35'
LS_COLORS+=':*.tiff=01;35'
LS_COLORS+=':*.png=01;35'
LS_COLORS+=':*.svg=01;35'
LS_COLORS+=':*.svgz=01;35'
LS_COLORS+=':*.mng=01;35'
LS_COLORS+=':*.pcx=01;35'
LS_COLORS+=':*.mov=01;35'
LS_COLORS+=':*.mpg=01;35'
LS_COLORS+=':*.mpeg=01;35'
LS_COLORS+=':*.m2v=01;35'
LS_COLORS+=':*.mkv=01;35'
LS_COLORS+=':*.webm=01;35'
LS_COLORS+=':*.ogm=01;35'
LS_COLORS+=':*.mp4=01;35'
LS_COLORS+=':*.m4v=01;35'
LS_COLORS+=':*.mp4v=01;35'
LS_COLORS+=':*.vob=01;35'
LS_COLORS+=':*.qt=01;35'
LS_COLORS+=':*.nuv=01;35'
LS_COLORS+=':*.wmv=01;35'
LS_COLORS+=':*.asf=01;35'
LS_COLORS+=':*.rm=01;35'
LS_COLORS+=':*.rmvb=01;35'
LS_COLORS+=':*.flc=01;35'
LS_COLORS+=':*.avi=01;35'
LS_COLORS+=':*.fli=01;35'
LS_COLORS+=':*.flv=01;35'
LS_COLORS+=':*.gl=01;35'
LS_COLORS+=':*.dl=01;35'
LS_COLORS+=':*.xcf=01;35'
LS_COLORS+=':*.xwd=01;35'
LS_COLORS+=':*.yuv=01;35'
LS_COLORS+=':*.cgm=01;35'
LS_COLORS+=':*.emf=01;35'
LS_COLORS+=':*.axv=01;35'
LS_COLORS+=':*.anx=01;35'
LS_COLORS+=':*.ogv=01;35'
LS_COLORS+=':*.ogx=01;35'
LS_COLORS+=':*.aac=00;36'
LS_COLORS+=':*.au=00;36'
LS_COLORS+=':*.flac=00;36'
LS_COLORS+=':*.m4a=00;36'
LS_COLORS+=':*.mid=00;36'
LS_COLORS+=':*.midi=00;36'
LS_COLORS+=':*.mka=00;36'
LS_COLORS+=':*.mp3=00;36'
LS_COLORS+=':*.mpc=00;36'
LS_COLORS+=':*.ogg=00;36'
LS_COLORS+=':*.ra=00;36'
LS_COLORS+=':*.wav=00;36'
LS_COLORS+=':*.axa=00;36'
LS_COLORS+=':*.oga=00;36'
LS_COLORS+=':*.spx=00;36'
LS_COLORS+=':*.xspf=00;36'
export LS_COLORS
