# profile.d/git-cinnabar.sh
export PATH="$PATH:$HOME/git/dotfiles/submodules/git-cinnabar"
