# profile.d/less.sh
export LESS=-iRFX

export LESS_TERMCAP_so=$(printf "\e[7m")
export LESS_TERMCAP_se=$(printf "\e[0m")
