# profile.d/pyenv.sh
export PATH="/home/calid/.pyenv/bin:$PATH"
eval "$(pyenv init -)"
eval "$(pyenv virtualenv-init -)"
