# profile.d/python.sh
export PYTHON_INCLUDE_PATH="/usr/include/python3.5/"
PATH=~/.local/bin:$PATH
